# script.conf
# simulating 30 different programms
[program:script1]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script2]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script3]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script4]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script5]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script6]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script7]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script8]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script9]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script10]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script11]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script12]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script13]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script14]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script15]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script16]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script17]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script18]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script19]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script20]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script21]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script22]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script23]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script24]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script20]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script21]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script22]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script23]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script24]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script25]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script26]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script27]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script28]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script29]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:script30]
command=python3 /var/app/script.py
directory=/var/app
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[group:script]
programs=script1,script2,script3,script4,script5,script6,script7,script8,script9,script10,script11,script12,script13,script14,script15,script16,script17,script18,script19,script20,script21,script22,script23,script24,script25,script26,script27,script28,script29,script30
